# DockerFiles/06-centos_httpd/conf/supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock 

[supervisord]
nodaemon=true

[program:httpd]
command=/usr/sbin/httpd -DFOREGROUND

[group:allservices]
programs=httpd

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock